# carhole_web.service
[Unit]
Description=Carhole Web Frontend
Requires=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/carhole_web
ExecStart=/usr/local/bin/bundle exec puma -b tcp://0.0.0.0:8080 -e production
TimeoutSec=30
RestartSec=15

[Install]
WantedBy=multi-user.target
